# ocr-dashboard.service
[Unit]
Description=OCR Web Dashboard
After=network.target postgresql.service

[Service]
User=tomaasz
Group=tomaasz
WorkingDirectory=/home/tomaasz/ocr-dashboard-v3

# Load environment variables from .env file
EnvironmentFile=-/home/tomaasz/ocr-dashboard-v3/.env

Environment="PATH=/home/tomaasz/ocr-dashboard-v3/venv/bin:/usr/bin"
Environment="PYTHONPATH=/home/tomaasz/ocr-dashboard-v3"
# Note: %% is used to escape % in systemd unit files
Environment="OCR_PG_DSN=postgresql://tomaasz:123Karinka!%%40%%23@127.0.0.1:5432/ocr"

# Standardized source mount root (all hosts mount storage here)
Environment="OCR_SOURCE_ROOT=/data/sources"
Environment="OCR_SOURCE_DIR=/mnt/nas_genealogy/Sources"

# Desktop SSH (Windows -> WSL)
Environment="OCR_REMOTE_DESKTOP_HOST=1dtw2.tail7319.ts.net"
Environment="OCR_REMOTE_DESKTOP_USER=tomaa"
Environment="OCR_REMOTE_DESKTOP_REPO_DIR=C:\\Dev\\ocr-dashboard-v3"
Environment="OCR_REMOTE_DESKTOP_SOURCE_DIR=W:\\Sources\\Nurskie dokumenty"
Environment="OCR_REMOTE_DESKTOP_WSL_DISTRO=Ubuntu-24.04"
Environment="OCR_WSL_SHARE_PROFILE=1"
Environment="OCR_REMOTE_DESKTOP_SSH_OPTS=-i /home/tomaasz/.ssh/ocr_desktop -o BatchMode=yes -o ConnectTimeout=8"

# WSL SSH (WSL has its own sshd, exposed on port 2222)
Environment="OCR_REMOTE_HOST=1dtw2.tail7319.ts.net"
Environment="OCR_REMOTE_USER=tomaasz"
Environment="OCR_REMOTE_REPO_DIR=/home/tomaasz/ocr-dashboard-v3"
Environment="OCR_REMOTE_SOURCE_DIR=/mnt/kosciesza_sources/"
Environment="OCR_REMOTE_SSH_OPTS=-p 2222"

ExecStart=/home/tomaasz/ocr-dashboard-v3/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 9090
Restart=always
RestartSec=5

# --- SEKCJA LIMITÓW ZASOBÓW ---
# Najniższy priorytet procesora (SSH zawsze będzie miało pierwszeństwo)
Nice=19
# Ograniczenie do 6 z 8 wątków (600%), aby zostawić 2 wolne dla SSH/OS
CPUQuota=600%
# Limit RAM na 24GB (zostawiamy 8GB dla systemu)
MemoryMax=24G
# ------------------------------

[Install]
WantedBy=multi-user.target
