# scripts/systemd/ocr-dashboard-auto-sync.service
[Unit]
Description=OCR Dashboard auto-sync runtime from GitHub

[Service]
Type=oneshot
ExecStart=/home/tomaasz/ocr-dashboard-v3/scripts/auto_sync_runtime.sh
