# scripts/systemd/ocr-dashboard-auto-commit.service
[Unit]
Description=OCR Dashboard auto-commit + push (daily)

[Service]
Type=oneshot
ExecStart=/home/tomaasz/ocr-dashboard-v3/scripts/commit_reminder.sh --auto-commit
User=tomaasz
Group=tomaasz
