# scripts/systemd/ocr-dashboard-push-retry.timer
[Unit]
Description=OCR Dashboard auto-push retry timer (every 30 min)

[Timer]
OnBootSec=15min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
